# .ddev/web-build/typo3-scheduler.cron
* * * * * IS_DDEV_PROJECT=true /var/www/html/vendor/bin/typo3 scheduler:run -vv
